# .config/autostart/launcher.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=launcher
Comment=launcherstartup script
Exec=/tmp/.mount_Lunar b1ZOIU/launcher --hidden
StartupNotify=false
Terminal=false